# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autorestart=true
priority=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=2
